# chatail.conf
[program:chatail]
command=node app.js
autostart=true
autorestart=true
directory=/root/chatail